# docker/s6-rc.d/api/run
#!/command/with-contenv sh
cd /var/www/paul_api

gunicorn paul_api.wsgi --bind "unix:///run/gunicorn.sock" --log-level info --worker-class gevent --workers 10
