# docker/s6-rc.d/qcluster/run
#!/command/with-contenv sh
cd /var/www/paul_api

python manage.py qcluster
